# packaging/systemd/rtparityd.service
[Unit]
Description=rtparityd prototype storage daemon
Documentation=man:rtparityd(8) man:rtparityd.conf(5)
After=network.target local-fs.target
Wants=local-fs.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/sbin/rtparityd \
    --socket=/var/run/rtparityd/rtparityd.sock \
    --pid-file=/var/run/rtparityd/rtparityd.pid \
    --config=/etc/rtparityd/rtparityd.conf \
    --journal-path=/var/lib/rtparityd/journal.bin \
    --metadata-path=/var/lib/rtparityd/metadata.bin

ExecStop=/bin/rm -f /var/run/rtparityd/rtparityd.sock /var/run/rtparityd/rtparityd.pid
ExecReload=/bin/kill -HUP $MAINPID

Restart=on-failure
RestartSec=5s
TimeoutStartSec=30s
TimeoutStopSec=30s

RuntimeDirectory=rtparityd
RuntimeDirectoryMode=0755

PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/rtparityd /var/log/rtparityd

StandardOutput=journal
StandardError=journal
SyslogIdentifier=rtparityd

[Install]
WantedBy=multi-user.target
